# etc/proxmox-backup-proxy.service.in
[Unit]
Description=Proxmox Backup API Proxy Server
Wants=network-online.target
After=network.target
Wants=proxmox-backup.service
After=proxmox-backup.service

[Service]
Type=notify
ExecStart=%LIBEXECDIR%/proxmox-backup/proxmox-backup-proxy
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/proxmox-backup/proxy.pid
Restart=on-failure
User=%PROXY_USER%
Group=%PROXY_USER%

[Install]
WantedBy=multi-user.target
